# etc/init/battleship.conf
# copy this file as root to /etc/init and put battleship in the directory
# referenced below

description "Battleship server"
author      "Zachary L. Stauber"

# used to be: start on startup
# until we found some mounts weren't ready yet while booting

start on runlevel [2345]
stop on runlevel [!2345]

# automatically respawn

respawn
respawn limit 99 5
umask 022

script
    
	export HOME="/root"
	exec /usr/local/bin/node /opt/battleship/server.js >> /var/log/battleship.log 2>&1

end script

post-start script
						   
	# optionally put a script here that will notifiy you node has (re)started
	# /root/bin/hoptoad.sh "node.js has started!"
										    
end script
